# flaskapp.conf
[program:flaskapp]
command = python /opt/app.py
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/app1.log
stderr_logfile = /var/log/supervisor/app1.log
